# system/etc/profile.d/environ.sh
# The environment variable whose name includes lower case letter(s) should be
# defined here.
# The environment variable whose name consist from upper case letters at all
# should be defined at System's Property -> Detail Settings -> Environment Variable
#export http_proxy; http_proxy=http://proxy:3128/
#export https_proxy; https_proxy=http://proxy:3128/
#export ftp_proxy; ftp_proxy=http://proxy:3128/
